# systemd/system/gsync-watch.service
[Unit]
Description=GSync watch service

[Service]
ExecStart=/home/gsync/watch-today.sh 
User=gsync


[Install]
WantedBy=multi-user.target
