# cron.d/gsync
0  0    * * *   root    service gsync-watch restart
